# components/moby-dick/src/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.3 -F
priority=20
